# deployment/systemd/marxist-search-api.service
# Systemd service file for Marxist Search API
# Copy to: /etc/systemd/system/marxist-search-api.service

[Unit]
Description=Marxist Search API Server
Documentation=https://github.com/domwxyz/marxist-search
After=network.target

[Service]
Type=simple
User=marxist
Group=marxist
WorkingDirectory=/opt/marxist-search/backend

# Environment
Environment="PATH=/opt/marxist-search/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/marxist-search/backend"
Environment="ENVIRONMENT=production"
Environment="DEBUG=false"
Environment="TRANSFORMERS_TRUST_REMOTE_CODE=1"

# Data paths (production)
Environment="DATA_DIR=/var/lib/marxist-search"
Environment="DATABASE_PATH=/var/lib/marxist-search/articles.db"
Environment="INDEX_PATH=/var/lib/marxist-search/txtai"
Environment="CACHE_PATH=/var/lib/marxist-search/cache"

# Service execution
ExecStart=/opt/marxist-search/venv/bin/uvicorn \
    src.api.main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 3 \
    --log-level info \
    --no-access-log

# Restart policy
Restart=always
RestartSec=5
StartLimitInterval=0

# Resource limits
LimitNOFILE=65536
MemoryMax=6G

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/marxist-search /var/log/marxist-search

# Logging
StandardOutput=append:/var/log/marxist-search/api.log
StandardError=append:/var/log/marxist-search/errors.log
SyslogIdentifier=marxist-search-api

[Install]
WantedBy=multi-user.target
